# crontab
# Run yt-transcribe pipeline daily at 8:00 AM
0 8 * * * root cd /app && ./process_videos.sh >> /var/log/cron.log 2>&1

# Empty line required at end of crontab
